# irods/files/profile.d/postgres.sh
PATH=/usr/pgsql-12/bin/:"$PATH"
export PATH
